# templates/automysqlbackup.cron.erb
#!/bin/sh

BIN_FILE=/usr/local/bin/automysqlbackup

if `test -x $BIN_FILE`; then
	$BIN_FILE
	<% if cron_set_permissions %>
	chown root.root /var/backup/db* -R
	find <%= backup_dir %>/* -type f -exec chmod 400 {} \;
	find <%= backup_dir %>/* -type d -exec chmod 700 {} \;
	<% end %>
else
	echo "ERROR: File does not exist or `whoami` does not have permission to execute it."
	exit 1
fi
